# devops/roles/icos.restic_server/templates/restic-server.socket
[Socket]
ListenStream = {{ restic_server_listen }}

[Install]
WantedBy = sockets.target
