# devops/roles/icos.pve_server/templates/icos-auto-dnat.path
[Path]
PathChanged={{ pve_dnat_leases }}

[Unit]
Description=Detect when dnsmasq leases change and sync firewall rules

[Install]
WantedBy=multi-user.target
