# devops/roles/icos.pve_server/templates/icos-auto-dnat.timer
[Unit]
Description=Automatic port-forwarding for proxmox VMs

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
